# systemd/sshfs-mount.service
# Bitrix CDN SSHFS Mount Service
# Author: Chibilyaev Alexandr <info@aachibilyaev.com>
# Company: AAChibilyaev LTD

[Unit]
Description=SSHFS mount for Bitrix files
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
Environment=BITRIX_SERVER_USER=www-data
Environment=BITRIX_SERVER_IP=192.168.1.10
Environment=BITRIX_UPLOAD_PATH=/var/www/bitrix/upload

# ИСПРАВЛЕНО: Команда монтирования с правильными переменными
ExecStart=/usr/local/bin/mount-bitrix.sh

# Команда размонтирования
ExecStop=/usr/bin/fusermount -u /mnt/bitrix

# Перезапуск при сбое
Restart=on-failure
RestartSec=10

# Проверка что mount работает
ExecStartPost=/bin/sh -c 'sleep 2 && ls /mnt/bitrix > /dev/null'

[Install]
WantedBy=multi-user.target
